# specs/default/chef/site-cookbooks/ocs/templates/default/sgemasterd.service.erb
[Unit]
Description=Univa Grid Engine Qmaster
After=network.target remote-fs.target

[Service]
Type=forking
ExecStart=/bin/sh <%= @ocsroot %>/<%= @ocscell %>/common/sgemaster start
ExecStop=/bin/sh <%= @ocsroot %>/<%= @ocscell %>/common/sgemaster stop

[Install]
WantedBy=multi-user.target
